# Storage/drill/1.22.0/24.03-lts-sp1/drill.ini
[program:drill]
command=sh -c '/apache-drill/bin/drillbit.sh start && sleep 2 && tail -f /dev/null /apache-drill/log/*'